# prometheus/rootfs/etc/profile.d/prometheus.sh
export PROMETHEUS_HOME=/opt/prometheus
export PATH=${PATH}:${PROMETHEUS_HOME}